# debian/postinst
#!/bin/sh
# postinst script for composer-global-update
#
# see: dh_installdeb(1)

#DEBHELPER#

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

# Source debconf library.

. /usr/share/debconf/confmodule

db_version 2.0

reloadwebserver() {
        if [ -f /etc/init.d/nginx ] ; then
            if [ -x /usr/sbin/invoke-rc.d ]; then
                invoke-rc.d nginx reload 3>/dev/null || true
            else
                /etc/init.d/nginx reload 3>/dev/null || true
            fi
        fi

        if [ -f /etc/init.d/lighttpd ] ; then
            if [ -x /usr/sbin/invoke-rc.d ]; then
                invoke-rc.d lighttpd reload 3>/dev/null || true
            else
                /etc/init.d/lighttpd reload 3>/dev/null || true
            fi
        fi

        if [ -f /etc/init.d/apache2 ] ; then
            if [ -x /usr/sbin/invoke-rc.d ]; then
                invoke-rc.d apache2 reload 3>/dev/null || true
            else
                /etc/init.d/apache2 reload 3>/dev/null || true
            fi
        fi
}


case "$1" in
    configure)

	# Update Apps?
	db_get composer/UPDATE_APPS || true
	UPDATE_APPS=$RET

	db_get composer/WEB_USER || true
	WEB_USER=$(echo $RET | awk -F: '{print $1}')
	WEB_GROUP=$(echo $RET | awk -F: '{print $2}')
	if [ -z ${WEB_GROUP} ]; then
	    WEB_GROUP=${WEB_USER}
	fi


    # creating www-data group if it isn't already there
    if ! getent group $WEB_GROUP >/dev/null; then
        addgroup --system --force-badname --quiet $WEB_GROUP
    fi

    # creating www-data user if it isn't already there
    if ! getent passwd $WEB_USER >/dev/null; then
        adduser --system --force-badname --quiet \
            --ingroup ${WEB_GROUP} \
            --home /var/www --no-create-home \
            --shell /usr/sbin/nologin \
            $WEB_USER
        usermod -c "Web data" $WEB_USER
    fi

    if [ -z $WEB_GROUP ] ; then
        chown $WEB_USER /var/lib/composer/ -R
    else
        chown $WEB_USER:$WEB_GROUP /var/lib/composer/ -R
    fi

    if [ "$UPDATE_APPS" = "true" ]; then
	composer-global-update
        reloadwebserver
    fi

    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

exit 0
