# debian/postinst
#!/bin/sh
# postinst script for composer-global-update
#
# see: dh_installdeb(1)

#DEBHELPER#

set -e

# summary of how this script can be called:
#        * <postinst> `configure' <most-recently-configured-version>
#        * <old-postinst> `abort-upgrade' <new version>
#        * <conflictor's-postinst> `abort-remove' `in-favour' <package>
#          <new-version>
#        * <postinst> `abort-remove'
#        * <deconfigured's-postinst> `abort-deconfigure' `in-favour'
#          <failed-install-package> <version> `removing'
#          <conflicting-package> <version>
# for details, see http://www.debian.org/doc/debian-policy/ or
# the debian-policy package

# Source debconf library.
. /usr/share/debconf/confmodule

reloadwebserver() {
        if [ -f /etc/init.d/nginx ] ; then
            if [ -x /usr/sbin/invoke-rc.d ]; then
                invoke-rc.d nginx reload 3>/dev/null || true
            else
                /etc/init.d/nginx reload 3>/dev/null || true
            fi
        fi

        if [ -f /etc/init.d/lighttpd ] ; then
            if [ -x /usr/sbin/invoke-rc.d ]; then
                invoke-rc.d lighttpd reload 3>/dev/null || true
            else
                /etc/init.d/lighttpd reload 3>/dev/null || true
            fi
        fi

        if [ -f /etc/init.d/apache2 ] ; then
            if [ -x /usr/sbin/invoke-rc.d ]; then
                invoke-rc.d apache2 reload 3>/dev/null || true
            else
                /etc/init.d/apache2 reload 3>/dev/null || true
            fi
        fi
}


case "$1" in
    configure)

    mkdir -p /var/lib/composer/

    db_input low composer/WEB_USER || true
    db_go || true

    db_get composer/WEB_USER
    chown $RET:$RET /var/lib/composer/ -Rv

    db_get composer/UPDATE_APPS
    if [ "$RET" = "true" ]; then
	composer-global-update
        reloadwebserver
    fi

    ;;

    abort-upgrade|abort-remove|abort-deconfigure)
    ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
    ;;
esac

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

exit 0
